# bckup.service
[Unit]
Description=run backup creator script

[Service]
Type=simple
ExecStart=~/bin/bckup.sh

[Install]
WantedBy=multi-user.target
